# dot_config/systemd/user/openclaw-gateway.service.d/env.conf
[Service]
EnvironmentFile=/home/rai/.openclaw/.env
